# systemd/ordr.fm.timer
[Unit]
Description=Run ordr.fm music organization every 6 hours
Requires=ordr.fm.service

[Timer]
# Run every 6 hours
OnCalendar=*-*-* 00,06,12,18:00:00
# Run 5 minutes after boot if we missed a scheduled run
OnBootSec=5min
# If the service fails, retry after 30 minutes
OnUnitInactiveSec=30min
# Randomize start time by up to 15 minutes to avoid system load spikes
RandomizedDelaySec=15min
# Persist across reboots
Persistent=true

[Install]
WantedBy=timers.target